# deploy/RunAll.service
[Unit]
Description=run-all

[Service]
Type=forking
ExecStart=/usr/bin/RunAll.sh

[Install]
WantedBy=multi-user.target